# mksplash-theme-change.path
[Unit]
Description=Watch plymouthd.conf for theme changes

[Path]
PathChanged=/etc/plymouth/plymouthd.conf

[Install]
WantedBy=multi-user.target
